# deployment/configuration/parqueryd.target
[Unit]
Description=Parqueryd Main Service

[Install]
WantedBy=multi-user.target
